# recipes-milestone/sleep-mode/sleep-mode/sleep_mode.service
[Unit]
Description=low power mode service

[Service]
ExecStart=/usr/sbin/sleep_mode/sleep_mode.sh

[Install]
WantedBy=multi-user.target